# supervisord/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --tries=2
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker_%(process_num)02d.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
[program:user-online]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan user:online
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/user_online_%(process_num)02d.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
